# server-confs.service
[Unit]
Description=Server confs service
After=v2ray.service

[Service]
ExecStart=/etc/server-confs.sh

[Install]
WantedBy=default.target
