# server-confs.service
[Unit]
Description=Server confs service

[Service]
ExecStart=/etc/server-confs.sh

[Install]
WantedBy=default.target
